# reverseqr.service
# ReverseQR systemd service file
# DEPRECATED: This file is deprecated. Use setup-service.sh instead.
# 
# Installation:
#   chmod +x setup-service.sh
#   sudo ./setup-service.sh
#
# This will generate a service file with the correct paths for your installation.
#
# Manual service management:
#   sudo systemctl status reverseqr
#   sudo systemctl restart reverseqr
#   sudo journalctl -u reverseqr -f

# NOTE: The below is a template. Use setup-service.sh to generate an actual service file.

[Unit]
Description=ReverseQR - Secure File Sharing Server
After=network.target

[Service]
Type=simple
User=armand
WorkingDirectory=/home/armand/Videos/reverseqr

# Environment variables (loaded from .env file)
EnvironmentFile=/home/armand/Videos/reverseqr/.env

# Start command
ExecStart=/usr/bin/node /home/armand/Videos/reverseqr/src/server.js

# Restart policy
Restart=on-failure
RestartSec=10

# Process management
StandardOutput=journal
StandardError=journal
SyslogIdentifier=reverseqr

# Security settings
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
